# roles/system/files/btrfs-balance@.timer
[Unit]
Description=Weekly Btrfs balance (light) on %I

[Timer]
OnCalendar=weekly
Persistent=true

[Install]
WantedBy=timers.target
