# roles/dev/templates/archdev-pull.timer.j2
[Unit]
Description=Timer for ArchDev GitOps Automator

[Timer]
OnBootSec=15min
OnUnitActiveSec=6h
Unit=archdev-pull.service

[Install]
WantedBy=timers.target
